# install/ai-receptionist-agent.conf
[program:ai-receptionist-agent]
directory=/var/www/fspbx/resources/ai-receptionist-agent
command=/opt/fspbx/ai-receptionist-agent/.venv/bin/python -B -m agent.openai_controller
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
user=www-data
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1"
redirect_stderr=true
stdout_logfile=/var/log/fspbx-ai-receptionist-agent.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
