# install/reverb.conf
[program:reverb]
process_name=%(program_name)s
command=php /var/www/fspbx/artisan reverb:start --host=127.0.0.1 --port=8095
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/fspbx/storage/logs/reverb.log
stopwaitsecs=10
